# BuildTools/crontab
@hourly source /var/local/bin/stocks_venv/bin/activate && python manage.py runjobs hourly
